# deploy/systemd/omni-council-weekly.timer
[Unit]
Description=Run Omni-BMO Council Weekly Report

[Timer]
OnCalendar=Sun *-*-* 16:00:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
